# config/etc/services.d/fluentd/run
#!/bin/sh

fluentd -c /etc/fluentd/fluent.conf
